# crontab
@reboot cd /app && python3 /app/start.py
*/10 * * * * cd /app && python3 /app/start.py
