# crontab
@reboot cd /app && python3 /app/main.py 
*/10 * * * *  cd /app && python3 /app/main.py 
