# crontab
@reboot cd /app && python3 /app/start_telegram.py 
*/10 * * * *  cd /app && python3 /app/start_telegram.py 
